# launchd/claude-transcripts.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{OB_LAUNCHD_PREFIX}}.claude-transcripts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/{{OB_PYTHON_VERSION}}/bin/python3</string>
        <string>{{OB_DIR}}/bin/process-transcripts.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{OB_DIR}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/{{OB_USERNAME}}</string>
        <key>PATH</key>
        <string>/Users/{{OB_USERNAME}}/.bun/bin:/Users/{{OB_USERNAME}}/.local/bin:/Users/{{OB_USERNAME}}/.nvm/versions/node/{{OB_NODE_VERSION}}/bin:/Library/Frameworks/Python.framework/Versions/{{OB_PYTHON_VERSION}}/bin:/Users/{{OB_USERNAME}}/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>OB_LOG_DIR</key>
        <string>{{OB_LOG_DIR}}</string>
        <key>OPENBRAIN_URL</key>
        <string>{{OPENBRAIN_URL}}</string>
        <key>OPENBRAIN_KEY</key>
        <string>{{OPENBRAIN_KEY}}</string>
    </dict>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>{{OB_LOG_DIR}}/transcripts.log</string>

    <key>StandardErrorPath</key>
    <string>{{OB_LOG_DIR}}/transcripts-error.log</string>
</dict>
</plist>
