# launchd/claude-dreaming.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{OB_LAUNCHD_PREFIX}}.claude-dreaming</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/{{OB_USERNAME}}/.local/bin/claude -p "$(cat {{OB_DIR}}/prompts/dreams-prompt.md)" --dangerously-skip-permissions --output-format json; {{OB_DIR}}/bin/dreaming-post.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{OB_DIR}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/{{OB_USERNAME}}</string>
        <key>PATH</key>
        <string>/Users/{{OB_USERNAME}}/.bun/bin:/Users/{{OB_USERNAME}}/.local/bin:/Users/{{OB_USERNAME}}/.nvm/versions/node/{{OB_NODE_VERSION}}/bin:/Library/Frameworks/Python.framework/Versions/{{OB_PYTHON_VERSION}}/bin:/Users/{{OB_USERNAME}}/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_PATH</key>
        <string>/Users/{{OB_USERNAME}}/.nvm/versions/node/{{OB_NODE_VERSION}}/lib/node_modules</string>
        <key>OB_LOG_DIR</key>
        <string>{{OB_LOG_DIR}}</string>
        <key>OB_MEMORY_FILE</key>
        <string>{{OB_MEMORY_FILE}}</string>
        <key>OB_STATS_DIR</key>
        <string>{{OB_STATS_DIR}}</string>
        <key>OPENBRAIN_URL</key>
        <string>{{OPENBRAIN_URL}}</string>
        <key>OPENBRAIN_KEY</key>
        <string>{{OPENBRAIN_KEY}}</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{{OB_LOG_DIR}}/dreaming.log</string>

    <key>StandardErrorPath</key>
    <string>{{OB_LOG_DIR}}/dreaming-error.log</string>
</dict>
</plist>
